# docker/.user/.config/systemd/user/gnome-settings-daemon.target
[Unit]
Description=GNOME Settings Daemon Components
After=gnome-session.service
Requires=gnome-session.service
PartOf=graphical-session.target

[Install]
WantedBy=default.target